# infrastructure/templates/30-cloud-init/cluster-mgmt/systemd/omni-rotate-sa-key.timer
[Unit]
Description=Bi-Weekly rotation of Omni Service Account Key

[Timer]
OnUnitActiveSec=2w
RandomizedDelaySec=1h
Persistent=true

[Install]
WantedBy=timers.target
